# emoncms/service-runner
#!/sbin/openrc-run

depend() {
  need redis
}

name=$RC_SVCNAME
command=python3 /var/www/emoncms/scripts/services/service-runner/service-runner.py
